# tornado.conf.default
[group:tornadoes]
programs=tornado-9000,tornado-9001,tornado-9002,tornado-9003,tornado-9004,tornado-9005,tornado-9006,tornado-9007

[program:tornado-9000]
command=python /home/work/server/exhibit/server/router.py --port=9000
directory=/home/work/server/exhibit/server/
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/home/work/server/exhibit/log/tornado.log
loglevel=info

[program:tornado-9001]
command=python /home/work/server/exhibit/server/router.py --port=9001
directory=/home/work/server/exhibit/server/
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/home/work/server/exhibit/log/tornado.log
loglevel=info

[program:tornado-9002]
command=python /home/work/server/exhibit/server/router.py --port=9002
directory=/home/work/server/exhibit/server/
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/home/work/server/exhibit/log/tornado.log
loglevel=info

[program:tornado-9003]
command=python /home/work/server/exhibit/server/router.py --port=9003
directory=/home/work/server/exhibit/server/
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/home/work/server/exhibit/log/tornado.log
loglevel=info

[program:tornado-9004]
command=python /home/work/server/exhibit/server/router.py --port=9004
directory=/home/work/server/exhibit/server/
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/home/work/server/exhibit/log/tornado.log
loglevel=info

[program:tornado-9005]
command=python /home/work/server/exhibit/server/router.py --port=9005
directory=/home/work/server/exhibit/server/
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/home/work/server/exhibit/log/tornado.log
loglevel=info

[program:tornado-9006]
command=python /home/work/server/exhibit/server/router.py --port=9006
directory=/home/work/server/exhibit/server/
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/home/work/server/exhibit/log/tornado.log
loglevel=info

[program:tornado-9007]
command=python /home/work/server/exhibit/server/router.py --port=9007
directory=/home/work/server/exhibit/server/
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/home/work/server/exhibit/log/tornado.log
loglevel=info
